# deployment/supervisor.conf
[program:opencoindata]
directory=/root/opencoindata
command=/root/opencoindata/bin/opencoindata collect -i 10s
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600