# deployment/supervisor.conf
[program:opencoindata]
directory=/root/opencoindata
command=go run -u main.go poll
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=GOPATH="/root/.go"